# visualizer/pistream-hdmi-watch.service
[Unit]
Description=Synchrofazotron: start/stop the visualizer based on HDMI status
After=multi-user.target

[Service]
ExecStart=/opt/pistream-visualizer/hdmi-watch.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
