# visualizer/pistream-visualizer.service
[Unit]
Description=Synchrofazotron audio visualizer (cava / glslViewer) on the HDMI console
# We take over tty1 — the login console comes back with: systemctl start getty@tty1
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
# viz-run.sh picks the engine per /opt/pistream-visualizer/engine (panel-managed)
ExecStart=/opt/pistream-visualizer/viz-run.sh
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Environment=TERM=linux
# cava refuses to start without HOME (ERR_HOMELESS), even with an explicit -p
Environment=HOME=/root
Restart=on-failure
RestartSec=3

# Started exclusively by pistream-hdmi-watch — no [Install] section on purpose.
