# web/bt-agent.service
[Unit]
Description=Bluetooth auto-accept agent (headless "Just Works" pairing)
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
# NoInputNoOutput => "Just Works" pairing, no confirmation on the Pi side
ExecStart=/usr/bin/bt-agent -c NoInputNoOutput
Restart=on-failure
RestartSec=3
# bt-agent ignores SIGTERM (bluez-tools quirk) — without this every restart/
# shutdown waits the full 90 s stop timeout before SIGKILL. SIGINT works.
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=bluetooth.target
