# web/bt-agent.service
[Unit]
Description=Bluetooth auto-accept agent (headless "Just Works" pairing)
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
# NoInputNoOutput => "Just Works" pairing, no confirmation on the Pi side
ExecStart=/usr/bin/bt-agent -c NoInputNoOutput
Restart=on-failure
RestartSec=3

[Install]
WantedBy=bluetooth.target
